# deploy/debuginfod-go.service
[Unit]
Description=debuginfod-go debuginfo HTTP server (Astra Linux, Ubuntu, RedOS, CentOS)
Documentation=https://github.com/RioTwWks/debuginfod-go
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debuginfod
Group=debuginfod
EnvironmentFile=-/etc/debuginfod-go/debuginfod-go.env
ExecStart=/usr/bin/debuginfod
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# Hardening (см. deploy/security/README.md)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/var/lib/debuginfod-go /var/cache/debuginfod-go
# PostgreSQL через TCP — дополнительных путей не нужно.
# Unix-socket PostgreSQL: раскомментируйте:
# BindReadOnlyPaths=/var/run/postgresql

[Install]
WantedBy=multi-user.target
