# deploy/debuginfod-go.service
[Unit]
Description=debuginfod-go debuginfo HTTP server (Astra Linux, Ubuntu, RedOS, CentOS)
Documentation=https://github.com/RioTwWks/debuginfod-go
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debuginfod
Group=debuginfod
EnvironmentFile=-/etc/debuginfod-go/debuginfod-go.env
ExecStart=/usr/bin/debuginfod
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# Hardening (опционально, раскомментируйте при необходимости)
# NoNewPrivileges=true
# ProtectSystem=strict
# ReadWritePaths=/var/lib/debuginfod-go /var/cache/debuginfod-go

[Install]
WantedBy=multi-user.target
